# deploy/lumo-catalog-sync.timer
[Unit]
Description=Run the NS.gifts catalog stock sync periodically

[Timer]
# First run shortly after boot, then every 15 minutes. Persistent catches up a
# missed run after downtime.
OnBootSec=5min
OnUnitActiveSec=15min
Persistent=true

[Install]
WantedBy=timers.target
